# usbredir-helper/data/dbus-org.freedesktop.usbredir1.service.in
[Unit]
Description=USBRedir system helper

[Service]
Type=dbus
BusName=org.freedesktop.usbredir1
ExecStart=@ExecStart@
